# dev/config/AP/hostapd.service
[Service]
ExecStart=/usr/sbin/hostapd /etc/hostapd/hostapd.conf
RestartSec=5
Restart=on-failure
